# roles/ambient-light/ambient-light.service
[Unit]
Description=Ambient Light
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/etc/ambient-light
User=root
ExecStart=/usr/bin/ambient-light

[Install]
WantedBy=multi-user.target
